# packaging/systemd/bitcoind-exporter.service
[Unit]
Description=Bitcoin node Prometheus exporter
Documentation=https://github.com/rossigee/bitcoind-exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bitcoind-exporter
Group=bitcoind-exporter
EnvironmentFile=-/etc/default/bitcoind-exporter
ExecStart=/usr/bin/bitcoind-exporter
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/bitcoind-exporter

[Install]
WantedBy=multi-user.target
